# deploy/systemd/running-coach-morning-brief.service.example
[Unit]
Description=Send Running Coach morning Telegram brief
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/abel/Documentos/Proyectos/personal-trainer
ExecStart=/home/abel/Documentos/Proyectos/personal-trainer/.venv/bin/python scripts/notifications/coach_messages.py send-morning-brief
Environment=PATH=/home/abel/.opencode/bin:/home/abel/Documentos/Proyectos/personal-trainer/.venv/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
